# services/redis/run
#!/bin/bash
exec /sbin/setuser abc redis-server